# gov_last.conf
[program:gov_last]
command=bash -c "cd /root/gov_last/ && python3 main.py"
autostart=true
autorestart=true
stderr_logfile=/var/log/gov_last.err.log
stdout_logfile=/var/log/gov_last.out.log